# config/supervisor.d/crond.ini
[program:crond]
command = crond -f -L 15
autostart=true
autorestart=true
