# config/supervisor.d/postconf.ini
[program:postconf]
autostart       = false
autorestart     = false
command         = /usr/bin/postconf.sh
stdout_logfile=/var/log/supervisord-postconf.log
stderr_logfile=/var/log/supervisord-postconf.log
startsecs       = 0
